# luxon-sangkwon-mcp.service
[Unit]
Description=Luxon Sangkwon MCP Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=mcpserver
Group=mcpserver
WorkingDirectory=/opt/luxon-sangkwon-mcp
ExecStart=/usr/bin/python3 server.py --transport streamable-http --port 8102
EnvironmentFile=/opt/luxon-sangkwon-mcp/.env
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
